# board/rpi_common/overlay/etc/udev/rules.d/80-net.rules
SUBSYSTEM=="net", ENV{DEVTYPE}=="wlan", ENV{ID_BUS}=="usb", NAME="wlusb"
SUBSYSTEM=="net", RUN+="/etc/rc.netif"
