# board/rpi4/overlay/etc/udev/rules.d/81-dwc.rules
SUBSYSTEM=="udc", DEVPATH=="/devices/platform/soc/fe980000.usb/udc/fe980000.usb", RUN+="/etc/rc.usbif"
